# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid ${USER} /usr/local/bin/go-cron "%(ENV_CRON_SCHEDULE)s" /bin/bash -c "python /${APP_HOME}/backup.py"
