# root/etc/services.d/docker/run
#!/usr/bin/with-contenv bash

exec /usr/local/bin/wrapdocker